# services/lapis/build/supervisor/conf.d/lapis.conf
[program:lapis]
command = lapis server production
directory = /app
user = root
autostart = true
