# services/openresty/build/supervisor/conf.d/app.conf
[program:app]
command = /usr/local/openresty/bin/openresty -p /app -c nginx.conf
directory = /app
user = root
autostart = true
autorestart = true
