# services/lumen/build/supervisor/conf.d/app.conf
[program:php]
command = /usr/local/sbin/php-fpm
user = root
autostart = true
autorestart = true
