# services/react/build/supervisor/conf.d/app.conf
[program:php]
command = php /app/app.php
user = root
autostart = true
autorestart = true
