# services/net_core/build/supervisor/conf.d/app.conf
[program:app]
command = dotnet run
directory = /app
user = app
autostart = true
autorestart = true
